# ubuntu/upstart/iitk-fortigate.conf
description     "IITK Fortigate Authenticate"
start on started tty1
task
exec /usr/sbin/iitk-fortigate > /tmp/fortigate.log &
