# ubuntu/upstart/ironport.conf
description     "IITK Ironport Authenticate"
start on started tty1

expect fork
respawn

exec /usr/sbin/daemonize /usr/sbin/iitk-ironport
