# ubuntu/iitk-ironport.conf
description     "IITK Ironport Authenticate"
start on started tty1
task
exec /usr/sbin/iitk-ironport > /tmp/ironport.log &
